# adds/supervisord-matrix.conf
[program:matrix]
command=/usr/bin/python -m synapse.app.homeserver --config-path /data/homeserver.yaml
