# adds/supervisord-turnserver.conf
[program:turnserver]
command=/usr/bin/turnserver -c /data/turnserver.conf
